# deploy/systemd/metacleaner-bot.service
[Unit]
Description=MetaCleaner Telegram Bot
After=network.target redis.service metacleaner-worker.service
Wants=redis.service

[Service]
Type=simple
WorkingDirectory=/opt/metacleaner_bot

# SECURITY FIX: Запуск от непривилегированного пользователя
User=metacleaner
Group=metacleaner

# SECURITY FIX: Изоляция и ограничения
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
ReadWritePaths=/opt/metacleaner_bot/temp /opt/metacleaner_bot/logs /opt/metacleaner_bot/metacleaner.db /opt/metacleaner_bot/metacleaner.db-shm /opt/metacleaner_bot/metacleaner.db-wal

# Ограничения ресурсов
MemoryLimit=1G
CPUQuota=100%
TasksMax=100

# Переменные окружения
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/metacleaner_bot/.env

ExecStart=/opt/metacleaner_bot/.venv/bin/python -m bot.main

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

[Install]
WantedBy=multi-user.target
