# deploy/systemd/metacleaner-admin.service
[Unit]
Description=MetaCleaner Admin (FastAPI/Uvicorn)
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
WorkingDirectory=/root/metacleaner_bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/metacleaner_bot/.venv/bin/python -m admin
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
